# barcode-telegram-polling.conf
[program:barcode-telegram-polling]
command=/opt/php81/bin/php /var/www/telega/data/www/telega.4natic.ru/artisan telebot:polling
autostart=true
autorestart=true
stopasgroup=true
killasgroup=true
user=telega
numprocs=1
redirect_stderr=true
stdout_logfile=/var/www/telega/data/www/telega.4natic.ru/polling.log
stopwaitsecs=3600
